# systemd-user/.config/systemd/user/obsidian-backup.service
[Unit]
Description=Obsidian Vaults Backup with GDrive

[Service]
Type=oneshot
ExecStart=/home/lucas/.local/bin/obsidian-bisync.sh

